# infra/launcher/control-plane-quick-tunnel.service
# /etc/systemd/system/control-plane-quick-tunnel.service
#
# TryCloudflare quick tunnel — no Cloudflare account, no domain required.
# Exposes http://localhost:8789 (the launcher) on a random
# https://<random>.trycloudflare.com URL.
#
# WARNING: the URL changes every time cloudflared restarts. Suitable for
# bootstrapping / smoke tests only. For production, replace with a named
# tunnel + Cloudflare Access (see docs/DEPLOYMENT.md §14).
#
# After (re)start, fetch the current URL with:
#   /opt/hermes-control-plane/quick-tunnel-url.sh

[Unit]
Description=control-plane-launcher quick tunnel (TryCloudflare, ephemeral URL)
After=network-online.target control-plane-launcher.service
Wants=network-online.target
Requires=control-plane-launcher.service

[Service]
Type=simple
User=hermes-cp
Group=hermes-cp
# Cloudflare .deb installs to /usr/bin (not /usr/local/bin)
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate --url http://localhost:8789
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
