# infra/launcher/hermes-launcher.service
# /etc/systemd/system/hermes-launcher.service
#
# Install:
#   sudo cp hermes-launcher.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now hermes-launcher
#
# Logs:
#   sudo journalctl -u hermes-launcher -f
#
# Restart on env change:
#   sudo systemctl restart hermes-launcher

[Unit]
Description=Hermes launcher sidecar (E2B SDK + GitHub App + user OAuth token)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hermes-cp
Group=hermes-cp
WorkingDirectory=/opt/hermes-control-plane

# All secrets and config live here. Mode 0600, owned by hermes-cp:hermes-cp.
# See infra/launcher/env.example for the full list.
EnvironmentFile=/etc/hermes-control-plane/launcher.env

# bun must be on PATH for the hermes-cp user. Default: ~/.bun/bin/bun.
ExecStart=/home/hermes-cp/.bun/bin/bun run /opt/hermes-control-plane/launcher.js

Restart=always
RestartSec=5

# Hardening — the launcher does not need a writable filesystem outside /tmp.
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
